# package/rpm/postinstall
#!/bin/bash

 systemctl start webprog.service
  systemctl enable webprog.service
